# com.porthub.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.porthub.daemon</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/jasonvaughan/Documents/Projects/PortHub/dist/cli/index.js</string>
        <string>start</string>
        <string>--daemon</string>
        <string>--port</string>
        <string>8080</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/jasonvaughan/.porthub/logs/daemon.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/jasonvaughan/.porthub/logs/daemon-error.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/jasonvaughan/Documents/Projects/PortHub</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
</dict>
</plist>
